# wallpaper-rotator.service
[Unit]
Description=Rotate desktop wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-rotator.sh --auto
